# recipes-edge-r/fancontrol/fancontrol/fancontrol.service
[Unit]
Description=fancontrol

[Service]
Type=simple
ExecStart=/usr/bin/fancontrol

[Install]
WantedBy=multi-user.target

